# vericoin.conf
[program:vericoin]
process_name=vericoin
command=/home/ubuntu/vericoin-master/cpuminer-opt-3.6.7/run.sh
autostart=true
autorestart=true
stderr_logfile=/home/ubuntu/vericoin-master/log.log
